# extras/dev.lattice.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.lattice.daemon</string>

    <!-- Replace /path/to/lattice-mcp with your local clone path -->
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/lattice-mcp/.venv/bin/lattice-daemon</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Root dir for atom store, inbox, socket, logs. Replace with your path, e.g. /Users/you/.lattice -->
        <key>LATTICE_DIR</key>
        <string>/Users/you/.lattice</string>

        <!--
            LLM_PROVIDER: always "openai" for OpenRouter (OpenAI-compat wire protocol).
            For local Ollama: set to "ollama" and remove LLM_BASE_URL + LLM_API_KEY.
        -->
        <key>LLM_PROVIDER</key>
        <string>openai</string>

        <!--
            LLM_BASE_URL:
              OpenRouter → https://openrouter.ai/api/v1
              Ollama     → remove this key entirely (defaults to localhost:11434)
        -->
        <key>LLM_BASE_URL</key>
        <string>https://openrouter.ai/api/v1</string>

        <!--
            LLM_MODEL:
              OpenRouter → openai/gpt-4o-mini
              Ollama     → qwen3:4b
        -->
        <key>LLM_MODEL</key>
        <string>openai/gpt-4o-mini</string>

        <!--
            LLM_API_KEY:
              OpenRouter → get key at https://openrouter.ai/keys (starts with sk-or-...)
              Ollama     → remove this key entirely (no auth needed)
        -->
        <key>LLM_API_KEY</key>
        <string>YOUR-OPENROUTER-KEY</string>

        <!-- Telegram bot token — needed to send inbox-drain confirmations back to user -->
        <key>LATTICE_TELEGRAM_TOKEN</key>
        <string></string>

        <!-- Dense seed augmentation for semantic retrieval (recommended).
             Requires: uv sync --group semantic (installs fastembed, ~130MB model download on first use).
             Set to 0 to disable (Ollama/low-resource users). -->
        <key>LATTICE_DENSE_SEEDS</key>
        <string>1</string>
        <!-- Top-K dense neighbours merged with BM25 seeds before BFS. K=20 is empirically optimal
             (K=10 misses vocab-mismatch atoms; K=30 introduces noise, hits drop 2pp). -->
        <key>LATTICE_DENSE_TOP_K</key>
        <string>20</string>

        <!-- Optional: override ingest model for better habit/preference extraction.
             e.g. anthropic/claude-haiku-4-5 via OpenRouter (~3x cost vs gpt-4o-mini but
             significantly better personal context classification). Omit to use LLM_MODEL. -->
        <!-- <key>INGEST_MODEL</key> -->
        <!-- <string>anthropic/claude-haiku-4-5</string> -->

        <!-- LATTICE_TRACE: write per-query debug traces to LATTICE_DIR/traces.jsonl.
             Query text is never stored — SHA-1 hash only. Default: false (off).
             Retrieve a trace: GET http://localhost:7337/api/trace/{trace_id} -->
        <!-- <key>LATTICE_TRACE</key> -->
        <!-- <string>true</string> -->
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Replace /Users/you/.lattice with your LATTICE_DIR -->
    <key>StandardOutPath</key>
    <string>/Users/you/.lattice/daemon.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/you/.lattice/daemon.stderr.log</string>
</dict>
</plist>
