# extras/dev.lattice.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.lattice.daemon</string>

    <!-- Replace /path/to/lattice-mcp with your local clone path -->
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/lattice-mcp/.venv/bin/lattice-daemon</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Root dir for atom store, inbox, socket, logs. Replace with your path, e.g. /Users/you/.lattice -->
        <key>LATTICE_DIR</key>
        <string>/Users/you/.lattice</string>

        <!--
            LLM_PROVIDER: always "openai" for OpenRouter (OpenAI-compat wire protocol).
            For local Ollama: set to "ollama" and remove LLM_BASE_URL + LLM_API_KEY.
        -->
        <key>LLM_PROVIDER</key>
        <string>openai</string>

        <!--
            LLM_BASE_URL:
              OpenRouter → https://openrouter.ai/api/v1
              Ollama     → remove this key entirely (defaults to localhost:11434)
        -->
        <key>LLM_BASE_URL</key>
        <string>https://openrouter.ai/api/v1</string>

        <!--
            LLM_MODEL:
              OpenRouter → openai/gpt-4o-mini
              Ollama     → qwen3:4b
        -->
        <key>LLM_MODEL</key>
        <string>openai/gpt-4o-mini</string>

        <!--
            LLM_API_KEY:
              OpenRouter → get key at https://openrouter.ai/keys (starts with sk-or-...)
              Ollama     → remove this key entirely (no auth needed)
        -->
        <key>LLM_API_KEY</key>
        <string>YOUR-OPENROUTER-KEY</string>

        <!-- Telegram bot token — needed to send inbox-drain confirmations back to user -->
        <key>LATTICE_TELEGRAM_TOKEN</key>
        <string></string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Replace /Users/you/.lattice with your LATTICE_DIR -->
    <key>StandardOutPath</key>
    <string>/Users/you/.lattice/daemon.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/you/.lattice/daemon.stderr.log</string>
</dict>
</plist>
